# neutrixBuster/autostart/plank.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Plank
Comment=plank
Exec=plank
Icon=plank
Path=/usr/bin
Terminal=false
StartupNotify=false
